# supervisord.conf
[supervisord]
nodaemon=true              ; run in foreground (required for Docker)
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:gunicorn-django]
command=bash -c "gunicorn secsnow.wsgi -b 0.0.0.0:8000 --worker-class gevent --workers ${GUNICORN_WORKERS:-3} --worker-connections ${GUNICORN_WORKER_CONNECTIONS:-500} --timeout ${GUNICORN_TIMEOUT:-200} --keep-alive 5 --max-requests 1000 --max-requests-jitter 100"
directory=/opt/cloud/secsnow
stdout_logfile=log/gunicorn.django.log
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=10
priority=100
# 日志文件权限
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
redirect_stderr=true
environment=GUNICORN_WORKERS="%(ENV_GUNICORN_WORKERS)s",GUNICORN_WORKER_CONNECTIONS="%(ENV_GUNICORN_WORKER_CONNECTIONS)s",GUNICORN_TIMEOUT="%(ENV_GUNICORN_TIMEOUT)s"

